# templates/systemd/openclaw.service.j2
[Unit]
Description=OpenClaw AI Assistant
Documentation=https://openclaw.ai/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ openclaw_user }}
Group={{ openclaw_group }}
WorkingDirectory={{ openclaw_upstream_dir }}
ExecStart={{ openclaw_cli_bin }} gateway --port {{ openclaw_gateway_port }}
Restart=on-failure
RestartSec=10
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
EnvironmentFile={{ openclaw_env_path }}

# Security hardening (équilibré - sécurise sans bloquer l'agent)
{% if openclaw_security_hardening %}
PrivateTmp=yes
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
{% endif %}

# Resource limits
LimitNOFILE={{ openclaw_nofile_limit }}
{% if openclaw_memory_max | length > 0 %}
MemoryMax={{ openclaw_memory_max }}
{% endif %}

[Install]
WantedBy=multi-user.target
